# roles/cucumber/files/profile.d/cucumber.sh
alias capybara-gen=/usr/local/cucumber-pi/capybara-gen
alias cucumber="cucumber -I classic"
